# deploy/polymarket-specialist.service
[Unit]
Description=Polymarket Specialist Edge Strategy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/polymarket-copytrading-bot
ExecStart=/root/polymarket-copytrading-bot/.venv/bin/python scripts/run_specialist_strategy.py --run
EnvironmentFile=/root/polymarket-copytrading-bot/.env
Environment=PYTHONPATH=/root/polymarket-copytrading-bot
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polymarket-specialist

[Install]
WantedBy=multi-user.target
