# deploy/polymarket-roadmap-updater.timer
[Unit]
Description=Daily roadmap snapshot update (06:00 UTC)

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
